# homelab/backup-sync/backup-sync.timer
[Unit]
Description=Daily SSD-to-SD fallback sync

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
